# .config/systemd/user/wallpaper.service
[Unit]
Description=Change desktop wallpaper
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=%h/bin/set-random-wallpaper
